# src/systemd/clevis-luks-askpass.service.in
[Unit]
Description=Clevis LUKS systemd-ask-password Responder

[Service]
Type=oneshot
ExecStart=@bindir@/clevis luks askpass
