# src/luks/systemd/clevis-luks-askpass.path
[Unit]
Description=Clevis systemd-ask-password Watcher
DefaultDependencies=no
Before=cryptsetup-pre.target
Wants=cryptsetup-pre.target

[Path]
PathChanged=/run/systemd/ask-password

[Install]
WantedBy=cryptsetup.target
