# src/luks/systemd/clevis-luks-askpass.service.in
[Unit]
Description=Forward Password Requests to Clevis
Documentation=man:clevis-luks-unlockers(7)
DefaultDependencies=no

[Service]
Type=simple
ExecStart=@libexecdir@/clevis-luks-askpass -l
