# src/luks/systemd/clevis-luks-askpass.path
[Unit]
Description=Clevis systemd-ask-password Watcher
DefaultDependencies=no
Conflicts=shutdown.target
Before=basic.target shutdown.target remote-fs-pre.target

[Path]
PathChanged=/run/systemd/ask-password

[Install]
WantedBy=remote-fs.target
