# src/luks/systemd/clevis-luks-askpass.service.in
[Unit]
Description=Clevis LUKS systemd-ask-password Responder
DefaultDependencies=no
Conflicts=shutdown.target
Before=systemd-ask-password-console.service systemd-ask-password-plymouth.service shutdown.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/clevis-luks-askpass -l
