# src/luks/systemd/clevis-luks-askpass.service.in
[Unit]
Description=Clevis LUKS systemd-ask-password Responder
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=@libexecdir@/clevis-luks-askpass -l
