# src/luks/systemd/clevis-luks-askpass.service.in
[Unit]
Description=Clevis LUKS systemd-ask-password Responder
Documentation=man:clevis-luks-unlockers(7)
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=@libexecdir@/clevis-luks-askpass -l
